# vagrant/bash/files/etc/profile.d/common_aliases.sh
# Allow `sudo` to run on following aliases!
#   Tip from https://askubuntu.com/questions/22037/aliases-not-available-when-using-sudo
alias sudo='sudo '

# Only allow ENGLISH man pages, as they are much up-to-date and complete, and
# serve as reference!
alias man='LANG=C man'

